# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/log/supervisor/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24
environment=DISPLAY=":99"
autorestart=true
priority=1
user=agent
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb-error.log

[program:fluxbox]
command=/usr/bin/fluxbox -display :99
environment=DISPLAY=":99"
autorestart=true
priority=2
user=agent
stdout_logfile=/var/log/supervisor/fluxbox.log
stderr_logfile=/var/log/supervisor/fluxbox-error.log

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -nopw -rfbport 5900
environment=DISPLAY=":99"
autorestart=true
priority=3
user=agent
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc-error.log

[program:websockify]
command=websockify --web /opt/noVNC 0.0.0.0:6080 localhost:5900
autorestart=true
priority=4
user=agent
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify-error.log

[program:uvicorn]
# Note: --reload flag is controlled by UVICORN_RELOAD env var
# Set UVICORN_RELOAD=true in docker-compose.dev-override.yml for hot reload
command=/bin/sh -c 'if [ "$UVICORN_RELOAD" = "true" ]; then uvicorn orchestrator.api.main:app --host 0.0.0.0 --port 8001 --reload --reload-dir /app/orchestrator; else uvicorn orchestrator.api.main:app --host 0.0.0.0 --port 8001; fi'
directory=/app
environment=DISPLAY=":99",PYTHONPATH="/app"
autorestart=true
priority=10
user=agent
stdout_logfile=/var/log/supervisor/uvicorn.log
stderr_logfile=/var/log/supervisor/uvicorn-error.log

[program:agent_worker]
# Agent worker runs Claude CLI tasks from Redis queue
# This process runs outside uvicorn's context to avoid subprocess I/O issues
# numprocs=2 allows 2 concurrent agent tasks (matches typical pipeline: plan + generate)
command=python -m orchestrator.services.agent_worker
directory=/app
environment=DISPLAY=":99",PYTHONPATH="/app",AGENT_WORKER_ID="worker-%(process_num)02d"
autorestart=true
priority=11
user=agent
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/agent_worker_%(process_num)02d.log
stderr_logfile=/var/log/supervisor/agent_worker_%(process_num)02d-error.log
# Start after uvicorn to ensure Redis connection is ready
startsecs=5
# Graceful shutdown timeout
stopwaitsecs=30
stopasgroup=true
killasgroup=true

[program:autonomous_mission_worker]
# Temporal worker for persistent autonomous testing missions
command=/bin/sh -c 'if [ "$AUTONOMOUS_MISSION_WORKER_ENABLED" = "true" ]; then exec python -m orchestrator.services.autonomous_mission_worker; else echo "Autonomous mission worker disabled"; exec sleep infinity; fi'
directory=/app
environment=DISPLAY=":99",PYTHONPATH="/app"
autorestart=true
priority=12
user=agent
stdout_logfile=/var/log/supervisor/autonomous_mission_worker.log
stderr_logfile=/var/log/supervisor/autonomous_mission_worker-error.log
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true

[program:custom_workflow_worker]
# Temporal worker for durable custom workflow runs
command=/bin/sh -c 'if [ "$CUSTOM_WORKFLOW_WORKER_ENABLED" = "false" ]; then echo "Custom workflow worker disabled"; exec sleep infinity; else exec python -m orchestrator.services.custom_workflow_worker; fi'
directory=/app
environment=DISPLAY=":99",PYTHONPATH="/app"
autorestart=true
priority=13
user=agent
stdout_logfile=/var/log/supervisor/custom_workflow_worker.log
stderr_logfile=/var/log/supervisor/custom_workflow_worker-error.log
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true

[program:custom_browser_workflow_worker]
# Temporal worker for browser-backed agent runs that must render in VNC.
command=/bin/sh -c 'if [ "$CUSTOM_WORKFLOW_WORKER_ENABLED" = "false" ]; then echo "Custom browser workflow worker disabled"; exec sleep infinity; else exec python -m orchestrator.services.custom_workflow_worker; fi'
directory=/app
environment=DISPLAY=":99",PYTHONPATH="/app",TEMPORAL_WORKFLOW_TASK_QUEUE="%(ENV_TEMPORAL_BROWSER_WORKFLOW_TASK_QUEUE)s"
autorestart=true
priority=14
user=agent
stdout_logfile=/var/log/supervisor/custom_browser_workflow_worker.log
stderr_logfile=/var/log/supervisor/custom_browser_workflow_worker-error.log
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
